# project_empros/tests/deploy/systemd/sentinel-ingress.service
[Unit]
Description=Sentinel Nexus Core Ingress
After=network-online.target nats.service
Wants=network-online.target

[Service]
Type=notify
User=sentinel
Group=sentinel
WorkingDirectory=/opt/sentinel-nexus/dist
EnvironmentFile=/opt/sentinel-nexus/nexus.env
Environment="MALLOC_CONF=background_thread:true,metadata_thp:auto"
ExecStart=/opt/sentinel-nexus/dist/core_ingress

# Auto-recovery
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Watchdog -- restart if unresponsive for 30s
WatchdogSec=30

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/sentinel-nexus/dist

[Install]
WantedBy=multi-user.target
